# whisparr/rootfs/etc/s6-overlay/s6-rc.d/whisparr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Whisparr
# Runs Whisparr
# ==============================================================================

bashio::log.info 'Starting Whisparr...'

exec /opt/bin/Whisparr -nobrowser -data=/config
